# systemd.service.example
# Пример systemd сервиса для автозапуска Gunicorn
# Скопируйте этот файл в /etc/systemd/system/nesivtoroi-tech.service
# ВАЖНО: Замените 'ваш_пользователь' на имя вашего пользователя на сервере!

[Unit]
Description=Gunicorn instance для nesivtoroi.tech Flask приложения
After=network.target

[Service]
User=ваш_пользователь
Group=www-data
WorkingDirectory=/var/www/nesivtoroi.tech
Environment="PATH=/var/www/nesivtoroi.tech/venv/bin"
ExecStart=/var/www/nesivtoroi.tech/venv/bin/gunicorn --config /var/www/nesivtoroi.tech/gunicorn_config.py web_app:app

# Автоматический перезапуск при сбое
Restart=always
RestartSec=10

# Логирование
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nesivtoroi-tech

# Безопасность
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/nesivtoroi.tech

[Install]
WantedBy=multi-user.target

